# com.tiipeng.iagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tiipeng.iagent</string>

    <!-- venv created by install.sh; the symlink inside venv points to whichever
         Python the installer found (3.9.x default on Procursus). -->
    <key>ProgramArguments</key>
    <array>
        <string>/var/jb/var/mobile/iagent/venv/bin/python</string>
        <string>/var/jb/var/mobile/iagent/code/main.py</string>
    </array>

    <!-- iOS LaunchDaemons in the system domain reject UserName=mobile with
         exit 78 (EX_CONFIG). The daemon therefore runs as root. Files in
         IAGENT_HOME are mobile-owned, but root can read+write them. -->

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/var/mobile</string>
        <key>IAGENT_HOME</key>
        <string>/var/jb/var/mobile/iagent</string>
        <!-- Required so Python can verify TLS certificates on iOS -->
        <key>SSL_CERT_FILE</key>
        <string>/var/jb/etc/ssl/cert.pem</string>
        <key>PATH</key>
        <string>/var/jb/usr/bin:/var/jb/bin:/usr/bin:/bin</string>
        <!-- Force unbuffered stdout/stderr so logs are flushed even when
             the process is SIGKILL'd by Jetsam under memory pressure. -->
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/var/jb/var/mobile/iagent</string>

    <key>StandardOutPath</key>
    <string>/var/jb/var/mobile/iagent/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/jb/var/mobile/iagent/logs/stderr.log</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically on crash. -->
    <key>KeepAlive</key>
    <true/>

    <!-- Throttle restarts. 60s gives Jetsam time to release memory before
         we try again, and reduces the chance of launchd disabling the
         service after a flurry of OOM kills. -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Allow a few extra seconds for graceful shutdown on unload. -->
    <key>ExitTimeOut</key>
    <integer>10</integer>
</dict>
</plist>
