# com.tiipeng.iagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tiipeng.iagent</string>

    <!-- launchd invokes /bin/sh (symlinked through to Procursus's sh on
         rootless dopamine) explicitly, with the wrapper script as its
         argument. Going through `sh script.sh` instead of relying on
         the script's shebang avoids EX_CONFIG (78) from launchd's
         stricter exec policy. The wrapper itself drops privileges to
         the mobile user before exec'ing python — root-launched daemons
         get SIGKILL'd by iOS AMFI when they touch the network. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>/var/jb/var/mobile/iagent/code/daemon_wrapper.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/var/mobile</string>
        <key>IAGENT_HOME</key>
        <string>/var/jb/var/mobile/iagent</string>
        <key>SSL_CERT_FILE</key>
        <string>/var/jb/etc/ssl/cert.pem</string>
        <key>PATH</key>
        <string>/var/jb/usr/bin:/var/jb/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/var/jb/var/mobile/iagent</string>

    <key>StandardOutPath</key>
    <string>/var/jb/var/mobile/iagent/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/jb/var/mobile/iagent/logs/stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>30</integer>
    <key>ExitTimeOut</key>
    <integer>60</integer>
</dict>
</plist>
