# com.tiipeng.iagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tiipeng.iagent</string>

    <!-- Single-shot tick: long-poll Telegram for ~25 s, process any
         messages, exit. launchd wakes us again every 30 s. iOS Jetsam
         can't kill what isn't running, so the bot survives indefinitely. -->
    <key>ProgramArguments</key>
    <array>
        <string>/var/jb/var/mobile/iagent/venv/bin/python</string>
        <string>/var/jb/var/mobile/iagent/code/tick.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/var/mobile</string>
        <key>IAGENT_HOME</key>
        <string>/var/jb/var/mobile/iagent</string>
        <key>SSL_CERT_FILE</key>
        <string>/var/jb/etc/ssl/cert.pem</string>
        <key>PATH</key>
        <string>/var/jb/usr/bin:/var/jb/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/var/jb/var/mobile/iagent</string>

    <key>StandardOutPath</key>
    <string>/var/jb/var/mobile/iagent/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/jb/var/mobile/iagent/logs/stderr.log</string>

    <!-- Run once at load (catches anything queued while daemon was down)
         and then every 30 s thereafter. With LONG_POLL_TIMEOUT=25 in
         tick.py, the bot wakes ~5 s before the next tick — so apparent
         response latency is near zero when active, capped at 30 s when
         idle. Crucially: the process exits between ticks, so RAM is
         released and Jetsam has nothing to hunt. -->
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>30</integer>

    <!-- If a tick takes longer than 60 s for any reason, kill it so the
         next interval can start cleanly. -->
    <key>ExitTimeOut</key>
    <integer>60</integer>
</dict>
</plist>
